# ipwatcher.service
[Unit]
Description=ipwatcher dynamic DNS updater for Cloudflare and AWS Route 53
Documentation=https://github.com/msyrus/ipwatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ipwatcher
Group=ipwatcher
WorkingDirectory=/opt/ipwatcher
ExecStart=/opt/ipwatcher/ipwatcher
Restart=on-failure
RestartSec=10
UMask=0077

# Environment
EnvironmentFile=-/opt/ipwatcher/.env
Environment="CONFIG_FILE=/opt/ipwatcher/config.yaml"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipwatcher

[Install]
WantedBy=multi-user.target
